# roles/zfs/templates/zfs-scrub.timer.j2
[Unit]
Description=Scrub ZFS pools weekly

[Timer]
OnCalendar=Wed 09:33

[Install]
WantedBy=multi-user.target


