# roles/zfs/templates/zfs-scrub.service.j2
[Unit]
Description=Scrub ZFS pools

[Service]
Type=oneshot
ExecStart=/usr/bin/zpool scrub {{ zfs_pools | join(" ") }}
